# nubis/puppet/files/shim.ini
[program:bugzilla-pulse-shim]
command=/usr/bin/python /opt/pulse/shims/bugzilla_simple_shim.py -r 5 /opt/pulse/config.ini
user=root
directory=/opt/pulse
numprocs=1
stdout_logfile=/var/log/bugzilla_simple_shim-supervisor.log
stderr_logfile=/var/log/bugzilla_simple_shim-supervisor.log
autostart=false
autorestart=true
startsecs=10
